# server/files/sctl-relay.service
[Unit]
Description=sctl Tunnel Relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/sctl supervise --config /etc/sctl/relay.toml
Restart=on-failure
RestartSec=5
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sctl-relay
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
